# yadaguru-api/supervisord.conf
[supervisord]
nodaemon=true

[program:yadaguru]
command=/entrypoint.sh
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
